# InjectorUploader/InjectorUploader/Package/Library/LaunchDaemons/tw.edu.mcu.cce.nrl.InjectorUploader.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>tw.edu.mcu.cce.nrl.InjectorUploader</string>
	<key>Program</key>
	<string>/Applications/Injector.app/InjectorUploader</string>
	<key>RunAtLoad</key>
	<true/>
	<key>UserName</key>
	<string>root</string>
	<key>GroupName</key>
	<string>wheel</string>
	<key>KeepAlive</key>
	<true/>
	<key>Disabled</key>
	<false/>
	<key>ThrottleInterval</key>
	<integer>3600</integer>
</dict>
</plist>
