# InjectorDaemon/InjectorDaemon/Package/Library/LaunchDaemons/tw.edu.mcu.cce.nrl.InjectorKillDaemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/bin/killall</string>
		<string>-9</string>
		<string>InjectorDaemon</string>
	</array>
	<key>Label</key>
	<string>tw.edu.mcu.cce.nrl.InjectorKillDaemon</string>
	<key>Program</key>
	<string>/usr/bin/killall</string>
	<key>RunAtLoad</key>
	<false/>
	<key>UserName</key>
	<string>root</string>
	<key>GroupName</key>
	<string>wheel</string>
	<key>LaunchOnlyOnce</key>
	<true/>
	<key>Disabled</key>
	<true/>
</dict>
</plist>
