# InjectorDaemon/InjectorDaemon/Package/Library/LaunchDaemons/tw.edu.mcu.cce.nrl.InjectorDaemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/bin/screen</string>
		<string>-d</string>
		<string>-m</string>
		<string>/Applications/Injector.app/InjectorDaemon</string>
	</array>
	<key>Label</key>
	<string>tw.edu.mcu.cce.nrl.InjectorDaemon</string>
	<key>Program</key>
	<string>/usr/bin/screen</string>
	<key>RunAtLoad</key>
	<true/>
	<key>UserName</key>
	<string>root</string>
	<key>GroupName</key>
	<string>wheel</string>
	<key>LaunchOnlyOnce</key>
	<true/>
	<key>Disabled</key>
	<false/>
</dict>
</plist>
